# .config/systemd/user/signatured.service
[Unit]
Description=Signature daemon

[Service]
Type=simple
ExecStart=/bin/sh -c '${HOME}/.local/bin/signatured'
